# supervisor.ini
[program:sse]
directory=/var/www/sse/
command=/var/www/sse/run.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/sse_stderr.log
stdout_logfile=/var/log/sse_stdout.log
